# infra/systemd/infrared-sensor.service
# InfraRed Sensor 컴포넌트 — systemd 유닛 파일
# 설계서 v7.0 §4: Agent 권한 분리 5컴포넌트 구조
#
# 역할: /proc 스캔, FIM, 프로세스 계보 수집 (루트 권한)
#   - CAP_SYS_PTRACE: /proc 읽기
#   - CAP_AUDIT_READ: auditd 이벤트 읽기
#   - 수집된 이벤트 → UDS로 collector 전달
#   - 네트워크 통신 없음
#
# 설치:
#   sudo cp infra/systemd/infrared-sensor.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now infrared-sensor

[Unit]
Description=InfraRed Sensor Component
Documentation=https://github.com/your-org/infrared
After=local-fs.target

[Service]
Type=simple
ExecStart=/opt/infrared/venv/bin/python -m infrared_agent.components.sensor
EnvironmentFile=/etc/infrared/agent.env
WorkingDirectory=/var/lib/infrared
RuntimeDirectory=infrared

Restart=on-failure
RestartSec=5s
StartLimitInterval=200
StartLimitBurst=5

# ── Cgroup 리소스 제한 ────────────────────────────────────────────────────────
CPUQuota=3%
MemoryMax=64M
MemorySwapMax=0
IOWeight=10

# ── 보안: 필요한 Capabilities만 허용 ─────────────────────────────────────────
User=root
# 최소 권한 — 나머지 불필요한 권한 제거
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_AUDIT_READ CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_SYS_PTRACE CAP_AUDIT_READ CAP_DAC_READ_SEARCH
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/infrared
ReadOnlyPaths=/proc /var/log
PrivateTmp=yes
PrivateDevices=yes

# 네트워크 없음 (UDS만 허용)
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=multi-user.target
